# deb/package/DEBIAN/prerm
#!/bin/bash
set -e

if systemctl is-active --quiet kloak 2>/dev/null; then
    systemctl stop kloak || true
fi
if systemctl is-enabled --quiet kloak 2>/dev/null; then
    systemctl disable kloak || true
fi
